# .config/systemd/user/pkg-manifest.service
[Unit]
Dscription=generate package manifests

[Service]
Type=oneshot
ExecStart=/home/bells/bin/pkg-manifest.sh
